# backend/systemd/chatagent.service
[Unit]
Description=ChatAgent — FastAPI server
After=network.target

[Service]
User=chatagent
Group=chatagent
WorkingDirectory=/home/chatagent/chatagent/src/chatagent/backend
EnvironmentFile=/home/chatagent/chatagent/src/chatagent/backend/.env
ExecStart=/home/chatagent/chatagent/src/chatagent/backend/.venv/bin/uvicorn chatagent.app:app --host ${HOST} --port ${PORT}
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
